# crontab
# m h  dom mon dow   command
* * * * * echo "cron running $(date)" >> /tmp/crontab.log
0 0 * * * ./manage.py find_duplicate_records
